# deploy/events_ingest.service
[Unit]
Description=Ingest events into newshub
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/default/mcp-news
WorkingDirectory=/opt/mcp-news
ExecStart=/opt/mcp-news/.venv/bin/python scripts/ingest_events.py
Restart=on-failure

[Install]
WantedBy=multi-user.target

