# deploy/mcp-news.service
[Unit]
Description=MCP News Server (stdio via systemd-run is typical)
After=network.target postgresql.service

[Service]
Type=simple
Environment=DATABASE_URL=postgresql://localhost/newshub
Environment=ENABLE_SERVER_EMBEDDING=1
Environment=EMBEDDING_MODEL=intfloat/multilingual-e5-base
Environment=EMBED_SPACE=bge-m3
WorkingDirectory=%h/mcp-news
ExecStart=%h/mcp-news/.venv/bin/python -m mcp_news.server
Restart=on-failure

[Install]
WantedBy=multi-user.target
