# docs/batch/deploy/embed.service
[Unit]
Description=Embed chunks into pgvector (oneshot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Per-user runtime dir to host a simple flock lock
RuntimeDirectory=mcp-news
Environment=DATABASE_URL=postgresql://localhost/newshub
Environment=EMBEDDING_MODEL=intfloat/multilingual-e5-base
Environment=EMBED_SPACE=bge-m3
WorkingDirectory=%h/mcp-news
ExecStartPre=/usr/bin/flock -n /run/mcp-news/embed.lock -c 'true'
ExecStart=%h/mcp-news/.venv/bin/python scripts/embed_chunks.py --space ${EMBED_SPACE} --batch 128
TimeoutStartSec=600

[Install]
WantedBy=default.target
